# conf/crontab
* * * * * docker exec bdchat-api-php php /var/www/api/artisan schedule:run > /var/log/artisan_cronlog.log 2>&1